# config/includes.chroot/etc/cron.d/asl-motd-issue
*/15 * * * * /usr/local/sbin/update-motd-issue.sh
