# config/includes.chroot/etc/cron.d/astdb
@daily root /usr/sbin/astdb.php
@reboot root sleep 60 ; /usr/sbin/astdb.php